# scripts/storage-sync/seafile-storage-sync.service
[Unit]
Description=Seafile Storage Migration Sync
Documentation=https://github.com/nicogits92/seafile-deploy
After=network-online.target
Wants=network-online.target

# Only start if migration is active
ConditionPathExists=/opt/seafile/.storage-migration.conf

[Service]
Type=simple
ExecStart=/opt/seafile/storage-sync.sh
Restart=on-failure
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=seafile-storage-sync

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
