# scripts/gitops/seafile-gitops-sync.service
[Unit]
Description=Seafile GitOps Webhook Listener — syncs .env from Gitea on push
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/opt/seafile/.env
ExecStart=/usr/bin/python3 /opt/seafile/seafile-gitops-sync.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
