# scripts/config-git-server/seafile-config-server.service
[Unit]
Description=Seafile Config Git Server — serves config history for Portainer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/seafile/seafile-config-server.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=seafile-config-server

[Install]
WantedBy=multi-user.target
